# packages/frabit-server.service
[Unit]
Description=frabitd
After=time-sync.target network.target

[Service]
EnvironmentFile=-/etc/sysconfig/frabit-server
Type=simple
User=frabit
Group=frabit
PermissionsStartOnly=true
ExecStart=/usr/bin/frabit-server start
ExecStop=/usr/bin/frabit-server stop
[Install]
WantedBy=multi-user.target
